# crontab
1 */3 * * * python -u /app/shuckscrape.py
